# supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=/loop.sh
autorestart=true
stdout_logfile=/var/log/main.log
stderr_logfile=/var/log/main_error.log
